# .docker/supervisord.conf
# supervisord unificado del monolito Neon Trails.
# Orquesta los procesos del contenedor; los logs van a stdout/stderr para que
# `docker logs` (y Dokploy) los capturen.
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:php-fpm]
command=/usr/sbin/php-fpm8.5 -F
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Servidor X virtual (TigerVNC) en :1 / puerto RFB 5901.
# SecurityTypes None + localhost: el VNC solo es accesible por websockify (local);
# el acceso externo es únicamente vía nginx /game/ (juego-arcade público sin prompt).
[program:xvnc]
command=/usr/bin/Xtigervnc :1 -geometry %(ENV_VNC_RESOLUTION)s -depth 24 -rfbport 5901 -SecurityTypes None -localhost yes -AlwaysShared -desktop NeonTrails
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Puente noVNC: expone el VNC como WebSocket en :6080 (nginx lo proxya en /game/).
[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5901
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Gestor de ventanas mínimo + el juego (java -jar). Ver .docker/vnc-startup.sh.
[program:game]
command=/usr/local/bin/start-game.sh
autorestart=true
priority=50
environment=DISPLAY=":1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
